# templates/archive.service
[Unit]
Decription=Restart archive

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl try-restart sjcapture@archive.service
