# templates/archive.timer
[Unit]
Description=restart archive every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
